# NeverNote/build/conf/supervisor.conf
[supervisord]
nodaemon=true

[program:nginx]
user=root 
command=/usr/sbin/nginx
autostart=true
autorestart=true

[program:php-fpm]
user=root
command=/usr/local/sbin/php-fpm -F
autostart=true
autorestart=true
